# src/cron/cronjob.cron
* * * * * cd /app/dist/cron && node cron_scrape_and_push.js >> /app/cron.log 2>&1